# pgcronjob.crontab
*/1 * * * * ~/src/pgcronjob/pgcronjob.sh
